# docker/supervisord/dym_app/dym-queues.conf
[program:dym-queues]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --timeout=360
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/queues.worker.log
